# .config/systemd/user/anonymize-faces.service
[Unit]
Description=Pixelize faces in images dropped into .pixelize/full or .pixelize/slim

[Service]
Type=oneshot
ExecStart=%h/scripts/pixelize-inbox.sh
StandardOutput=journal
StandardError=journal
