# .config/systemd/user/video-capturas.service
[Unit]
Description=Capture frames from videos dropped into .capturas/full, medium or low

[Service]
Type=simple
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/home/tcabrera/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/scripts/capturas-inbox.sh
StandardOutput=journal
StandardError=journal
Restart=on-abnormal
RestartSec=5
